# pdf-viewer-kiosk/nitrobox-api.service
[Unit]
Description=Nitrobox PDF API (FastAPI/Uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mobilecharging
WorkingDirectory=/home/mobilecharging/demo-mobile-charging/pdf-viewer-kiosk
Environment=DOWNLOAD_DIR=/home/mobilecharging/demo-mobile-charging/pdf-viewer-kiosk/downloads
Environment=STATE_DIR=/home/mobilecharging/demo-mobile-charging/pdf-viewer-kiosk/.state
Environment=DEFAULT_WAIT_SECONDS=120
Environment=DEFAULT_POLL_SECONDS=5
ExecStart=/usr/bin/env bash -lc 'cd /home/mobilecharging/demo-mobile-charging/pdf-viewer-kiosk && . .venv/bin/activate && uvicorn app:app --host 0.0.0.0 --port 8000'
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
